# etc/profile.d/amdgpu-multigpu.sh
# Unpinned multi-GPU: do not force DRI_PRIME at login.
# Apps may set DRI_PRIME=1|2 for render offload on headless cards.